# ubuntu/ledmath.conf
#
# This task runs at boot to start LEDmath

description     "start LEDmath"

start on startup

task
exec /var/lib/cloud9/LEDscape/debian/pruss_magic_jr
exec /var/lib/cloud9/LEDscape/bin/ledmath
